# deploy/kronos-trend.service
# systemd unit for the Kronos crypto trend bot (Binance testnet).
# Install on a VPS (run these FROM INSIDE the repo dir):
#   REPO=$(pwd)
#   sudo cp deploy/kronos-trend.service /etc/systemd/system/
#   sudo sed -i "s#CHANGE_REPO_PATH#$REPO#g; s/CHANGE_USER/$USER/g" /etc/systemd/system/kronos-trend.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now kronos-trend
#   journalctl -u kronos-trend -f          # follow logs (stdout -> journal)
#
# Output goes to the systemd journal; the bot also writes rotating file logs to
# <repo>/logs/ and an hourly equity CSV.
[Unit]
Description=Kronos crypto trend bot (Binance testnet)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=20

[Service]
Type=simple
User=CHANGE_USER
WorkingDirectory=CHANGE_REPO_PATH
Environment=PYTHONPATH=src:.
ExecStart=CHANGE_REPO_PATH/.venv/bin/python -m kronos_mt5.live.run_trend_binance
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
